# scripts/templates/com.fno.watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.fno.watch.{{PROJECT}}</string>
  <key>ProgramArguments</key>
  <array>
    <string>{{FNO_WATCH_PATH}}</string>
    <string>{{PROJECT}}</string>
    <string>{{REPO_ROOT}}</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>{{HOME}}</string>
  </dict>
  <key>StandardOutPath</key>
  <string>{{REPO_ROOT}}/.fno/abi-watch.out.log</string>
  <key>StandardErrorPath</key>
  <string>{{REPO_ROOT}}/.fno/abi-watch.err.log</string>
</dict>
</plist>
