# scripts/post-merge/com.fno.postmerge.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Per-repo post-merge watcher LaunchAgent (ab-4e9fb05a). Rendered by
  install.sh, which substitutes the double-brace placeholders for THIS repo.
  RunAtLoad is false and the installer never loads the agent: a human reviews
  the rendered plist and runs `launchctl load` themselves (the human gate).
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>{{LABEL}}</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>{{WATCH_SCRIPT}}</string>
  </array>

  <!-- launchd starts /bin/bash as a non-login, non-interactive shell with a
       minimal PATH (/usr/bin:/bin:/usr/sbin:/sbin), so bare `gh`/`jq`/`claude`
       lookups would fail and the job would exit before processing any merge.
       install.sh renders the INSTALL-TIME PATH here, which is where those tools
       were resolvable (Codex P1, PR #390). -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>{{PATH}}</string>
    <!-- Model for the per-merge `claude --print` ritual fire. Defaults to Haiku
         (cheap; the ritual is mechanical). Edit to `sonnet` for stronger triage
         judgment, or remove this key to inherit the CLI default. -->
    <key>POST_MERGE_MODEL</key>
    <string>{{MODEL}}</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>{{REPO_ROOT}}</string>

  <!-- Poll cadence in seconds. Default 600 (10 min); design suggests 5-15 min,
       cache-window-aware. Override at install time with POST_MERGE_INTERVAL. -->
  <key>StartInterval</key>
  <integer>{{INTERVAL}}</integer>

  <!-- Do NOT fire a poll the instant the agent loads; wait for the first
       interval. Keeps `launchctl load` side-effect-free until the next tick. -->
  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>{{LOG_OUT}}</string>
  <key>StandardErrorPath</key>
  <string>{{LOG_ERR}}</string>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
